# observer/systemd/observer.service
[Unit]
Description=Network Observer Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/network-test
ExecStart=/usr/bin/python3 -m observer.app
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1
Environment=OBSERVER_CONFIG=/home/pi/network-test/observer/config.yaml

[Install]
WantedBy=multi-user.target
